# mssql-supervisor-dev/supervisord.conf
[supervisord]
nodaemon=true

[program:mssql]
command=/opt/mssql/bin/sqlservr.sh